# systemd/vps-control-backup-weekly.timer
[Unit]
Description=Run AI-First VPS Control Layer Weekly Backup

[Timer]
OnCalendar=Sun *-*-* 04:15:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
